# Dockerfiles/ros-melodic-mac-vnc/conf.d/websockify.conf
[program:websockify]
command=websockify --web /usr/share/novnc 2020 localhost:5900
autorestart=true
